# os/openSUSE/root/bin/check_aulogs.service
[Unit]
Description=Check auditd logs

[Service]
ExecStart=/root/bin/check_aulogs.sh

[Install]
WantedBy=multi-user.target
